# app/Resources/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:php5-fpm]
command=/usr/sbin/php5-fpm

[program:postgres]
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
autorestart=true
user=postgres

[program:nginx]
command=/usr/sbin/nginx
